# dotfiles/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/services/battery-warning.service
[Unit]
Description=Battery Level Check Script

[Service]
Type=oneshot
ExecStart=%h/.config/caelestia/scripts/power-alerts.sh

