# dotfiles/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/caelestia/timers/battery-warning.timer
[Unit]
Description=Battery Warning Timer

[Timer]
OnBootSec=10sec
OnUnitActiveSec=10sec

[Install]
WantedBy=default.target
